# etc/systemd/system/power-battery.target
[Unit]
Description=System is running on battery
Conflicts=power-ac.target shutdown.target
