# etc/systemd/system/power-ac.target
[Unit]
Description=System is running on AC
Conflicts=power-battery.target shutdown.target
